# systemd_custom_unit/roles/custom-unit/files/watchlog.service
[Unit]
Description= this is a test unit file

[Service]
Type=oneshot
EnvironmentFile=/etc/sysconfig/watchlog
ExecStart=/opt/watchlog.sh $WORD $LOG
